# scripts/public-devnet-v1/testnet-frontend.service
[Unit]
Description=Permawrite experimental testnet frontend (Next.js)
After=network.target observer-rpc-proxy.service faucet-http.service
Wants=observer-rpc-proxy.service faucet-http.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/permawrite/testnet-frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=MFND_RPC_PROXY_UPSTREAM=http://127.0.0.1:8787/rpc
Environment=MFND_FAUCET_UPSTREAM=http://127.0.0.1:8788/faucet
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target